# live-build/config/includes.chroot/etc/systemd/system/twina-kiosk.service
[Unit]
Description=TwinaOS Kiosk Session
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=user
Group=user
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_SESSION_TYPE=wayland
Environment=WAYLAND_DISPLAY=wayland-0
ExecStartPre=/bin/mkdir -p /run/user/1000
ExecStartPre=/bin/chown user:user /run/user/1000
ExecStart=/usr/bin/cage -- /usr/bin/midori
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
